# deploy/transitops.service
[Unit]
Description=TransitOps Node.js Application
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/transitops/backend
ExecStart=__NPM_BIN__ start
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
